# systemd/tunnel-client.service
[Unit]
Description=内网穿透隧道客户端
Documentation=https://github.com/your-username/my_cloudflared
After=network.target nss-lookup.target
Wants=network.target

[Service]
Type=simple
User=tunnel
Group=tunnel
WorkingDirectory=/home/tunnel/my_cloudflared
ExecStart=/usr/bin/python3 /home/tunnel/my_cloudflared/client.py --server 144.202.26.208 --server-port 8000 --local 127.0.0.1 --local-port 5008 --subdomain p
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tunnel-client

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/tunnel/my_cloudflared

# 资源限制
LimitNOFILE=1024
LimitNPROC=512
MemoryMax=512M

# 环境变量
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/home/tunnel/my_cloudflared

[Install]
WantedBy=multi-user.target